# packages/amberelec/system.d/amberelec.target
[Unit]
Description=AmberELEC
Requires=multi-user.target graphical.target
After=graphical.target
Conflicts=rescue.target
AllowIsolate=yes

[Install]
Alias=default.target
