# packages/amberelec/system.d/hotkey.service
[Unit]
Description=Global hotkey

[Service]
ExecStart=/usr/bin/hotkey.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
