# packages/amberelec/system.d/lastgame.service
[Unit]
Description=Hotkey and save last game played in RetroArch

[Service]
ExecStart=/usr/bin/lastgame.sh
Restart=always
Type=idle

[Install]
WantedBy=multi-user.target