# projects/Rockchip/devices/RG351P/packages/odroidgoa-utils/system.d/battery.service
[Unit]
Description=RG351P Low Power monitor
ConditionPathExists=/sys/class/power_supply/battery/capacity
Before=amberelec.target

[Service]
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/battery.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
