# projects/Rockchip/devices/RG351P/packages/odroidgoa-utils/system.d/volume.service
[Unit]
Description=Volume button service
ConditionPathExists=/dev/input/event2
Before=351elec.target

[Service]
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/volume_sense.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
