# packages/ui/webui/system.d/webui.service
[Unit]
Description=AmberELEC Web UI
ConditionPathExists=|/storage/.cache/services/webui.conf

[Service]
Restart=on-failure
Environment=PYTHONPATH=/usr/share/webui
WorkingDirectory=/usr/share/webui
ExecStart=python3 /usr/share/webui/server.py

[Install]
WantedBy=multi-user.target
