# projects/Rockchip/devices/RG351MP/packages/odroidgoa-utils/system.d/adckeys.service
[Unit]
Description=External service

[Service]
ExecStart=/usr/bin/adckeys.sh
Restart=always
Type=idle

[Install]
WantedBy=multi-user.target