# packages/ui/351elec-webui/system.d/webui.service
[Unit]
Description=351ELEC Web UI
ConditionPathExists=|/storage/.cache/services/webui.conf

[Service]
Restart=on-failure
Environment=PYTHONPATH=/usr/share/351elec-webui
WorkingDirectory=/usr/share/351elec-webui
ExecStart=python3 /usr/share/351elec-webui/server.py

[Install]
WantedBy=multi-user.target
