# packages/351elec/system.d/fan_control.service
[Unit]
Description=Fan control

[Service]
ExecStart=/usr/bin/fan_control.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target 
