# packages/351elec/system.d/jslisten.target
[Unit]
Description=jslisten
Requires=multi-user.target
Conflicts=rescue.target
AllowIsolate=yes

[Install]
Alias=default.target
